# deploy/com.focusmate.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.focusmate.backend</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/kyleparker/Documents/New project 2/data/AI assistant/scripts/run_backend.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/kyleparker/Documents/New project 2/data/AI assistant/backend</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/kyleparker/Documents/New project 2/data/AI assistant/backend/data/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/kyleparker/Documents/New project 2/data/AI assistant/backend/data/launchd.err.log</string>
</dict>
</plist>
